# roles/jenkins_agents/templates/jenkins_agent.service.j2
[Unit]
Description=Jenkins Agent

[Service]
User={{ jenkins_user }}
WorkingDirectory={{ jenkins_agent_directory }}
ExecStart=java -jar "{{ jenkins_agent_directory }}/agent.jar" -jnlpUrl "http://{{ jenkins_master_url }}/computer/{{ inventory_hostname }}/jenkins-agent.jnlp"
Restart=always

[Install]
WantedBy=multi-user.target
